# Library/LaunchAgents/emacs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>emacs</string>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>Aqua</string>
		<string>Background</string>
		<string>LoginWindow</string>
		<string>StandardIO</string>
		<string>System</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
	  <key>COLORTERM</key>
	  <string>truecolor</string>
	  <key>LANG</key>
	  <string>en_US.UTF-8</string>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/opt/emacs/bin/emacs</string>
		<string>--fg-daemon</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
